# deploy/telemetry2.service.template
[Unit]
Description=Modbus Decoder / CG Telemetry
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User={{APP_USER}}
Group={{APP_USER}}
WorkingDirectory={{APP_DIR}}
ExecStart=/usr/bin/python3 app.py --config config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
